# skel/gui/.config/bash/profile.d/sway.bash
#!/bin/bash
export WALLPAPER="$HOME/media/wallpapers/fish.jpg"

## Sway TTY startup inserted from evolve/setup/gui.sh 
