# skel/gui/.config/bash/profile.d/sway.bash
#!/bin/bash
export WALLPAPER="$HOME/media/wallpapers/fish.jpg"

# used in sway config
export XDG_CONFIG_HOME="$HOME/.config"
## Sway TTY startup inserted from evolve/setup/gui.sh 
